# srcs/grafana/supervisor.conf
[supervisord]
nodaemon=true
user=root

[program:grafana]
command=sh -c "/grafana/bin/grafana-server"
autostart=true
autorestart=true
